# sysfiles/crontab
# m h dom mon dow user  command
1 0		* * *	bot		cd /home/bot/HSBot2/ && ./termine.py
8 23	* * 1	root	reboot
* *		* * *	bot		cd /homr/bot/HSBot2/ && ./dht11.py